# mpd.socket
[Socket]
ListenStream=127.0.0.1:6600

[Install]
WantedBy=default.target
